# systemd/reflection.service
# Magic Hat — Reflection AI Platform Service
# Copyright (c) 2026 George Scott Foley — MIT License

[Unit]
Description=Reflection AI Platform
Documentation=https://github.com/omegcrash/reflection
After=network-online.target postgresql.service redis.service
Wants=network-online.target
Requires=postgresql.service redis.service

[Service]
Type=notify
User=magichat
Group=magichat
WorkingDirectory=/home/magichat

# Environment
Environment=REFLECTION_ENV=production
Environment=REFLECTION_HOST=127.0.0.1
Environment=REFLECTION_PORT=8000
Environment=REFLECTION_WORKERS=4
EnvironmentFile=-/etc/magichat/reflection.env
EnvironmentFile=-/etc/magichat/providers.env

# Exec
ExecStartPre=/usr/bin/python3.11 -m reflection migrate
ExecStart=/usr/bin/python3.11 -m reflection serve \
    --host ${REFLECTION_HOST} \
    --port ${REFLECTION_PORT} \
    --workers ${REFLECTION_WORKERS}

# Watchdog
WatchdogSec=30
NotifyAccess=all

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=5

# Security hardening — CIS + OWASP Agentic Top 10
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/magichat/.familiar /home/magichat/.reflection /var/log/magichat
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
ProcSubset=pid
RestrictRealtime=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
MemoryDenyWriteExecute=yes
LockPersonality=yes
RemoveIPC=yes
UMask=0077

# System call filtering (OWASP ASI02 — tool sandboxing)
SystemCallFilter=@system-service
SystemCallFilter=~@mount @reboot @swap @clock @module @raw-io @debug
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM

# seccomp profile
# SeccompProfile=/etc/magichat/seccomp-reflection.json

# Network isolation — only loopback + outbound (no binding new ports)
IPAddressAllow=127.0.0.0/8 ::1/128
RestrictNetworkInterfaces=lo

# Capabilities — drop everything, grant only what's needed
CapabilityBoundingSet=
AmbientCapabilities=

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
LimitCORE=0
MemoryMax=2G
CPUQuota=80%
TasksMax=512

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=reflection

[Install]
WantedBy=multi-user.target
